# deploy/systemd/media-automata-deep-check.timer
[Unit]
Description=Run Media Automata deep Instagram native check hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=2min
Unit=media-automata-deep-check.service

[Install]
WantedBy=timers.target
