# deploy/systemd/media-automata-deep-check.service
[Unit]
Description=Media Automata Deep Instagram Native Check
After=network-online.target media-automata-api.service media-automata-redroid.service
Wants=network-online.target media-automata-redroid.service

[Service]
Type=oneshot
WorkingDirectory=/home/unichronic/media_automata
EnvironmentFile=/home/unichronic/media_automata/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/unichronic/media_automata/.venv/bin/python -m media_automata.cli monitor-once --deep-instagram
