# deploy/systemd/media-automata-monitor.timer
[Unit]
Description=Run Media Automata production monitor every two minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=2min
AccuracySec=15s
Unit=media-automata-monitor.service

[Install]
WantedBy=timers.target
